# .local/etc/zsh/profile.d/00-keychain.sh
#eval `keychain --agents ssh,gpg --dir "$XDG_DATA_HOME/keychain" --noask --quiet --eval <insert names of ssh or gpg keys here>`
