# .local/etc/zsh/profile.d/00-keychain.sh
#eval `keychain --agents ssh,gpg --eval <insert names of ssh or gpg keys here>`
